# Raspberry Pi/BHT-Wetter/Python/Python_cron/crontab
@reboot sleep 10; /usr/local/bin/python /app/crontab.py > /proc/1/fd/1 2>/proc/1/fd/2
